# recipes/tailscale/files/config/tailscale/var-lib-tailscale.mount
[Mount]
What=/config/tailscale/state
Where=/var/lib/tailscale
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
