# recipes/bird/files/config/bird/bird-watch.service
[Unit]
Description=Watch for File Changes and Restart BIRD

[Service]
ExecStart=/config/bird/bird-watch.sh

[Install]
WantedBy=multi-user.target
